# docker/legacy-portal/supervisor/conf.d/aggregators.conf
; AUTO-GENERATED by tools/gen-supervisor-fragments.cjs from docker-compose.yml.
; Aggregator-role programs (master + hybrid).
; Do not edit by hand — re-run the generator after changing compose backend services.

[program:alert-manager-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-alert-manager-aggregator.lck && exec sdcc-alert-manager-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asset-advertisement-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-asset-advertisement-aggregator.lck && exec sdcc-asset-advertisement-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asset-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-asset-aggregator.lck && exec sdcc-asset-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asset-health-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-asset-health-aggregator.lck && exec sdcc-asset-health-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:attack-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-attack-aggregator.lck && exec sdcc-attack-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-over-gre-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-over-gre-aggregator.lck && exec sdcc-bgp-over-gre-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:daily-jobs-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-daily-jobs-aggregator.lck && exec sdcc-daily-jobs-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:diversion-period-checker-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-diversion-period-checker-aggregator.lck && exec sdcc-diversion-period-checker-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:domain-resolve-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-domain-resolve-aggregator.lck && exec sdcc-domain-resolve-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-fail-over-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-fail-over-aggregator.lck && exec sdcc-dp-fail-over-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-policy-inconsistency-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-policy-inconsistency-aggregator.lck && exec sdcc-dp-policy-inconsistency-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-policy-manager]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-policy-manager.lck && exec sdcc-dp-policy-manager.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:file-reporter-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-file-reporter-aggregator.lck && exec sdcc-file-reporter-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:incident-manager]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-incident-manager.lck && exec sdcc-incident-manager.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pdf-incident-reporter-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-pdf-incident-reporter-aggregator.lck && exec sdcc-pdf-incident-reporter-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:resource-utilization-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-resource-utilization-aggregator.lck && exec sdcc-resource-utilization-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sc-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-sc-aggregator.lck && exec sdcc-sc-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:site-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-site-aggregator.lck && exec sdcc-site-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:static-route-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-static-route-aggregator.lck && exec sdcc-static-route-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:status-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-status-aggregator.lck && exec sdcc-status-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ts-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-ts-aggregator.lck && exec sdcc-ts-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-connectivity-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-connectivity-aggregator.lck && exec sdcc-vision-connectivity-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-cpe-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-cpe-aggregator.lck && exec sdcc-vision-cpe-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-device-policy-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-device-policy-aggregator.lck && exec sdcc-vision-device-policy-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-policies-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-policies-aggregator.lck && exec sdcc-vision-policies-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:waf-ssl-protection-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-waf-ssl-protection-aggregator.lck && exec sdcc-waf-ssl-protection-aggregator.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
