# docker/legacy-portal/supervisor/conf.d/minimal-hybrid.programs.conf
; AUTO-GENERATED by tools/gen-supervisor-fragments.cjs from docker-compose.yml.
; Minimal-hybrid set (pollers+aggregators); autostart only: cmd-executor, incident-manager.
; Do not edit by hand — re-run the generator after changing compose backend services.

[program:ar-config-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-ar-config-poller.lck && exec sdcc-ar-config-gaps.py"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:attacks-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-attacks-poller.lck && exec sdcc-attacks-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-advertise-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-advertise-poller.lck && exec sdcc-bgp-advertise-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-peer-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-peer-poller.lck && exec sdcc-bgp-peer-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-routing-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-route-poller.lck && exec sdcc-bgp-route-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cmd-executor]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-cmd-executor.lck && exec sdcc-cmd-executor.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:domain-resolve-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-domain-resolve-poller.lck && exec sdcc-domain-resolve-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-fail-over-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-fail-over-poller.lck && exec sdcc-dp-fail-over-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:flow-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-flow-poller.lck && exec sdcc-flow-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:health-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-health-poller.lck && exec sdcc-health-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:netflow-manager]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-netflow-manager.lck && exec sdcc-netflow-manager"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sc-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-sc-poller.lck && exec sdcc-sc-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:site-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-site-poller.lck && exec sdcc-site-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:static-routes-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-static-routes-poller.lck && exec sdcc-static-routes-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-api-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-api-poller.lck && exec sdcc-vision-api-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-control-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-control-poller.lck && exec sdcc-vision-control-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-device-policy-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-device-policy-poller.lck && exec sdcc-vision-device-policy-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-policies-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-policies-poller.lck && exec sdcc-vision-policies-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:waf-route-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-waf-route-poller.lck && exec sdcc-waf-route-poller.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:alert-manager-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-alert-manager-aggregator.lck && exec sdcc-alert-manager-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asset-advertisement-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-asset-advertisement-aggregator.lck && exec sdcc-asset-advertisement-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asset-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-asset-aggregator.lck && exec sdcc-asset-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asset-health-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-asset-health-aggregator.lck && exec sdcc-asset-health-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:attack-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-attack-aggregator.lck && exec sdcc-attack-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-over-gre-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-over-gre-aggregator.lck && exec sdcc-bgp-over-gre-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:daily-jobs-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-daily-jobs-aggregator.lck && exec sdcc-daily-jobs-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:diversion-period-checker-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-diversion-period-checker-aggregator.lck && exec sdcc-diversion-period-checker-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:domain-resolve-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-domain-resolve-aggregator.lck && exec sdcc-domain-resolve-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-fail-over-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-fail-over-aggregator.lck && exec sdcc-dp-fail-over-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-policy-inconsistency-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-policy-inconsistency-aggregator.lck && exec sdcc-dp-policy-inconsistency-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-policy-manager]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-policy-manager.lck && exec sdcc-dp-policy-manager.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:file-reporter-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-file-reporter-aggregator.lck && exec sdcc-file-reporter-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:incident-manager]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-incident-manager.lck && exec sdcc-incident-manager.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pdf-incident-reporter-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-pdf-incident-reporter-aggregator.lck && exec sdcc-pdf-incident-reporter-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:resource-utilization-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-resource-utilization-aggregator.lck && exec sdcc-resource-utilization-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sc-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-sc-aggregator.lck && exec sdcc-sc-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:site-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-site-aggregator.lck && exec sdcc-site-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:static-route-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-static-route-aggregator.lck && exec sdcc-static-route-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:status-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-status-aggregator.lck && exec sdcc-status-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ts-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-ts-aggregator.lck && exec sdcc-ts-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-connectivity-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-connectivity-aggregator.lck && exec sdcc-vision-connectivity-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-cpe-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-cpe-aggregator.lck && exec sdcc-vision-cpe-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-device-policy-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-device-policy-aggregator.lck && exec sdcc-vision-device-policy-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-policies-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-policies-aggregator.lck && exec sdcc-vision-policies-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:waf-ssl-protection-aggregator]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-waf-ssl-protection-aggregator.lck && exec sdcc-waf-ssl-protection-aggregator.py --debug true"
autostart=false
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
