# docker/legacy-portal/supervisor/conf.d/pollers.conf
; AUTO-GENERATED by tools/gen-supervisor-fragments.cjs from docker-compose.yml.
; Poller-role programs (monitor + hybrid).
; Do not edit by hand — re-run the generator after changing compose backend services.

[program:ar-config-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-ar-config-poller.lck && exec sdcc-ar-config-gaps.py"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:attacks-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-attacks-poller.lck && exec sdcc-attacks-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-advertise-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-advertise-poller.lck && exec sdcc-bgp-advertise-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-peer-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-peer-poller.lck && exec sdcc-bgp-peer-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bgp-routing-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-bgp-route-poller.lck && exec sdcc-bgp-route-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cmd-executor]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-cmd-executor.lck && exec sdcc-cmd-executor.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:domain-resolve-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-domain-resolve-poller.lck && exec sdcc-domain-resolve-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dp-fail-over-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-dp-fail-over-poller.lck && exec sdcc-dp-fail-over-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:flow-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-flow-poller.lck && exec sdcc-flow-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:health-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-health-poller.lck && exec sdcc-health-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:netflow-manager]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-netflow-manager.lck && exec sdcc-netflow-manager"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sc-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-sc-poller.lck && exec sdcc-sc-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:site-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-site-poller.lck && exec sdcc-site-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:static-routes-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-static-routes-poller.lck && exec sdcc-static-routes-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-api-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-api-poller.lck && exec sdcc-vision-api-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-control-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-control-poller.lck && exec sdcc-vision-control-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-device-policy-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-device-policy-poller.lck && exec sdcc-vision-device-policy-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vision-policies-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-vision-policies-poller.lck && exec sdcc-vision-policies-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:waf-route-poller]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-waf-route-poller.lck && exec sdcc-waf-route-poller.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
