# docker/legacy-portal/supervisor/minimal-hybrid.conf
; Supervisord config for an SDCC backend container.
;
; One config can run on any backend role: each sdcc-*.py app self-selects via
; @require_backend_roles against this container's SDCC_BACKEND_ROLE. Role-mismatched
; apps exit 0 immediately, so autorestart=unexpected keeps them down while genuinely
; crashed apps (non-zero exit) are restarted.
;
; Item 1 scope: cmd-executor + incident-manager only (role=hybrid covers both).
; More apps get added in later plan items.

[supervisord]
nodaemon=true
logfile=/var/log/sdcc/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

; Control socket so `supervisorctl` (and the `sdcc` shim that delegates to it) can
; query/start/stop programs.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:cmd-executor]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-cmd-executor.lck && exec sdcc-cmd-executor.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:incident-manager]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-incident-manager.lck && exec sdcc-incident-manager.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
