# docker/legacy-portal/supervisor/minimal-monitor.conf
; Supervisord profile: MINIMAL MONITOR — cmd-executor only, for a monitor backend.
;
; Mirrors minimal-hybrid but without incident-manager (which is an aggregator and
; requires master/hybrid role). cmd-executor is a poller-role program, so it passes
; @require_backend_roles on a monitor backend (SDCC_BACKEND_ROLE=monitor) and acts
; only on the SCs whose ScrubbingCenter.backend points at this backend.
;
; Use full-monitor.conf instead to run the complete poller set.

[supervisord]
nodaemon=true
logfile=/var/log/sdcc/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:cmd-executor]
command=/bin/sh -c "rm -f /var/lock/sdcc/sdcc-cmd-executor.lck && exec sdcc-cmd-executor.py --debug true"
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
